# docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)

[supervisord]
stdout_logfile=/dev/stdout   ; (main log file;default $CWD/supervisord.log)
stdout_logfile_maxbytes=0    ; Disabled to avoid rotation on /dev/stdout, which is illegal
stderr_logfile=/dev/stderr   ; (main log file;default $CWD/supervisord.log)
stderr_logfile_maxbytes=0    ; Disabled to avoid rotation on /dev/stdout, which is illegal
logfile_backups=0            ; Disabled to avoid rotation on /dev/stdout, which is illegal
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user=root		             ;

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:java]
command=java -jar /app/back/coddy-back.jar
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
